# bing-wallpaper.timer
[Unit]
Description="Auto Download Bing wallpaper Timer"

[Timer]
OnUnitActiveSec=1h
Unit=bing-wallpaper.service

[Install]
WantedBy = timers.target

